# osticket/supervisord.conf
[supervisord]
nodaemon=true

[program:php5-fpm]
command=php5-fpm -c /etc/php5/fpm

[program:nginx]
command=nginx

[eventlistener:cron]
command=php5 /data/upload/api/cron.php
events=TICK_60
